# resources/wayle.service
[Unit]
Description=Wayle Desktop Shell
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/wayle shell
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=30
StartLimitBurst=5

[Install]
WantedBy=graphical-session.target
